# deploy/cognitive-deploy.service
[Unit]
Description=Cognitive Core auto-deploy: git fetch + conditional reload
Documentation=https://github.com/mocartlex-wq/cognitive-core/blob/main/AGENT_OPERATIONS.md
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=salex
WorkingDirectory=/opt/cognitive-core
Environment=COGNITIVE_REPO_DIR=/opt/cognitive-core
Environment=COGNITIVE_DEPLOY_BRANCH=main

# Optional: Telegram alerts on deploy-failure / rollback.
# Create /etc/cognitive-deploy.env with TELEGRAM_BOT_TOKEN=... and TELEGRAM_CHAT_ID=...
# (chmod 600 root) — auto-deploy will pick them up. Without this file, alerts are
# silent (only journalctl).
EnvironmentFile=-/etc/cognitive-deploy.env

# Vision provider keys — owner-managed. Systemd reads as root *before* dropping
# to salex, then exports as env vars; docker compose interpolates them into
# the api container without ever touching the file directly (which would
# fail EACCES since the file is mode 600 root:root). See PR #57.
EnvironmentFile=-/etc/cogcore-qwen.env

ExecStart=/opt/cognitive-core/scripts/auto-deploy.sh

# Allow docker exec without password (sudoers configured separately)
# StandardOutput/Error → journal, читать через `journalctl -u cognitive-deploy -n 100 -f`
StandardOutput=journal
StandardError=journal

# Не дольше 10 минут — крайний случай для full rebuild
TimeoutStartSec=600
