# deploy/cognitive-orchestrator.service
[Unit]
Description=Cognitive Core multi-agent orchestrator (C2 MVP)
After=network-online.target cognitive-rooms.service
Wants=network-online.target

[Service]
Type=simple
User=salex
WorkingDirectory=/opt/cognitive-core
# /etc/cognitive-orchestrator.env должен задавать ORCH_ROOM_ID и ORCH_ROOM_KEY (секрет, не в git)
EnvironmentFile=/etc/cognitive-orchestrator.env
ExecStart=/usr/bin/python3 /opt/cognitive-core/scripts/cognitive-orchestrator.py --room ${ORCH_ROOM_ID} --watch 20
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
